# init/upstart/perflock.conf
# perflock - Benchmark performance locking service

description "perflock daemon"

start on runlevel [2345]
stop on runlevel [!2345]
respawn

console log

exec perflock -daemon
